# deploy/systemd/project-harpy-eagle-sync-results.timer
[Unit]
Description=Periodic Project Harpy Eagle S3 result sync

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Unit=project-harpy-eagle-sync-results.service

[Install]
WantedBy=timers.target

