# deploy/systemd/project-harpy-eagle-sync-results.service
[Unit]
Description=Sync Project Harpy Eagle results from S3
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=www-data
WorkingDirectory=/opt/project-harpy-eagle
EnvironmentFile=/etc/project-harpy-eagle-sync.env
ExecStart=/opt/project-harpy-eagle/scripts/sync_results_from_s3.sh

